# ops/saferide-alert.timer
[Unit]
Description=SafeRide Nepal hourly SecurityEvent alert check

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
