# ops/saferide-backup.service
[Unit]
Description=SafeRide Nepal database backup
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/saferide-backup-db
User=root
Group=root
